# docker/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=http://usuario:senha@proxy.example.com:80/"
